# wordpress-backup.timer
[Unit]
Description=Wordpress backup timer

[Timer]
OnCalendar=daily
Persistent=true
#OnUnitActiveSec=1s

[Install]
WantedBy=timers.target
